# gunicorn.conf
[program:gunicorn]
directory=/home/ubuntu/django/src/StudySmart
command=/home/ubuntu/env/bin/gunicorn --workers 3 --bind unix:/home/ubuntu/django/src/StudySmart/app.sock StudySmart.wsgi:application
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn/gunicorn.err.log
stdout_logfile=/var/log/gunicorn/gunicorn.out.log

[group:guni]
programs:gunicorn

